# trader/bitstar_api_bot.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

setuid kamil
setgid kamil

exec /home/kamil/.virtualenvs/bitstar-api-bot/bin/python /home/kamil/Documents/Bitstar/bitstar-api-bot/app.py
